# deb/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable proxmox-cpu-affinity
    systemctl restart proxmox-cpu-affinity
fi
